# etc/systemd/tlspool.socket
[Socket]
ListenStream=/var/run/tlspool.sock
SocketUser=tlspool
SocketGroup=tlspool
SocketMode=0666

[Install]
WantedBy=sockets.target
