# etc/systemd/tlspool.service
[Service]
ExecStart=/usr/local/bin/tlspool -c /usr/local/etc/tlspool.conf
